# supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=warn
logfile_maxbytes=0
logfile=/dev/stderr
pidfile=/tmp/supervisor.pid

[program:fpm]
priority=5
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
priority=10
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

